# linux/etc/system/ld-relay.service
[Unit]
Description=LaunchDarkly Relay Proxy

[Service]
Type=simple
ExecStart=/usr/bin/ld-relay

[Install]
WantedBy=multi-user.target
